# docs/bw-agent.service
[Unit]
Description=Bitwarden Desktop Agent
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/bw-agent
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
